# systemd/user/devops-full-health.service
[Unit]
Description=DevOps Playground - Full Health Check (oneshot)

[Service]
Type=oneshot
WorkingDirectory=%h/devops-playground
ExecStart=%h/devops-playground/full_health.sh

# Optional: keep output in journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
