# systemd/user/devops-full-health.service
[Unit]
Description=DevOps Playground - Full Health Check (oneshot)

[Service]
Type=oneshot
WorkingDirectory=%h/devops-playground

EnvironmentFile=%h/.config/devops-playground/devops-playground.env

ExecStart=%h/devops-playground/bin/dp health --quiet

# Optional: keep output in journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
